# services/cg-tcpdump.service
[Unit]
After=network.target

[Service]
Type=one-shot
ExecStart=/bin/capaguard/initTcpdump.sh

[Install]
WantedBy=default.target
